# crontab
0 23 * * * cd /app && /usr/local/bin/python clockify_automate_script.py >> /tmp/file 2>&1
